# multipathd/multipathd.socket
[Socket]
ListenStream=@/org/kernel/linux/storage/multipathd

[Install]
WantedBy=sockets.target
